# deploy/supervisord.d/collector.ini
[program:collector]
command=python3 /opt/collector/src/main.py
stdout_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stdout_logfile=/var/log/collector/collector.log
redirect_stderr=true
user=ingestion
priority=500
autostart=true
startsecs=10
startretries=100000
autorestart=true
exitcodes=0
